# root/etc/services.d/heimdalljs/run
#!/usr/bin/with-contenv bash

cd /app/heimdall || exit

exec \
    s6-setuidgid abc npm run server
