# root/etc/services.d/cloud9/run
#!/usr/bin/with-contenv bash

cd /c9sdk
HOME=/c9bins exec \
	s6-setuidgid abc \
	/c9bins/.c9/node/bin/node server.js --listen 0.0.0.0 -p 8000 -w /var/www/localhost/heimdall -a :
