# root/etc/s6-overlay/s6-rc.d/svc-queue/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc php /app/www/artisan queue:work database --sleep=3 --tries=3
